# docker/app/supervisord.conf
[supervisord]
nodaemon=true

[program:rails]
command=bin/rails s -b 0.0.0.0
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
